# port_trigger.service
# ============================================================
# 部署步骤（在 OrangePi 上执行）
#
# 1. 复制到 systemd 目录：
#      sudo cp port_trigger.service /etc/systemd/system/
#
# 2. 创建日志目录：
#      mkdir -p /home/orangepi/logs
#
# 3. 启用开机自启并立即启动：
#      sudo systemctl daemon-reload
#      sudo systemctl enable port_trigger.service
#      sudo systemctl start  port_trigger.service
#
# 4. 查看运行状态：
#      sudo systemctl status port_trigger
#
# 5. 实时查看日志：
#      tail -f /home/orangepi/logs/port_trigger.log
#
# 6. 停止 / 禁用自启：
#      sudo systemctl stop    port_trigger.service
#      sudo systemctl disable port_trigger.service
# ============================================================

[Unit]
Description=ROS Port Trigger Service
After=network.target

[Service]
Type=simple
User=orangepi
WorkingDirectory=/home/orangepi/slam_ws
ExecStart=/usr/bin/python3 /home/orangepi/slam_ws/port_trigger.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/orangepi/logs/port_trigger.log
StandardError=append:/home/orangepi/logs/port_trigger.log

[Install]
WantedBy=multi-user.target
